# scripts/trellis-restarter.service
[Unit]
Description=Trellis Restarter — Watches for restart trigger and restarts Ivy
Documentation=https://github.com/kymorrand/trellis
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kyle
Group=kyle
WorkingDirectory=/home/kyle/projects/trellis
EnvironmentFile=/home/kyle/projects/trellis/.env
ExecStart=/home/kyle/projects/trellis/scripts/trellis_restarter.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=trellis-restarter

[Install]
WantedBy=multi-user.target
